# rpi-setup-scripts/switch-to-client.service
[Unit]
Description=Switch to WiFi Client Mode on Startup
After=network-pre.target
Before=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/switch-to-client.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
